# data/io.github.pureblue_os.purebazaar.service.in
[Unit]
Description=Bazaar Service
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=io.github.pureblue_os.purebazaar
ExecStart=@bindir@/bazaar --no-window
